# physlock.service
[Unit]
Description=Physlock on Suspend

[Service]
User=bsuth
ExecStart=/usr/bin/physlock

[Install]
WantedBy=sleep.target
