# slock.service
[Service]
User=bsuth
ExecStart=/usr/bin/slock

[Install]
WantedBy=sleep.target
